# scicat-dataset-ingestor.service
[Unit]
Description=SciCat Dataset Ingerstor
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/jkotan/ndts/scingestor/scicat_dataset_ingestor

[Install]
WantedBy=multi-user.target
